# app-services/shogun-core/shogun-core.service
[Unit]
Description=Shogun Core — Japan trip AI concierge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=devops-agent
WorkingDirectory=/opt/git/work/shogun/app-services/shogun-core
EnvironmentFile=/opt/git/work/shogun/app-services/shogun-core/.env
ExecStart=/opt/git/work/shogun/app-services/shogun-core/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8082 --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shogun-core

[Install]
WantedBy=multi-user.target
